# v2g-liberty/rootfs/etc/s6-overlay/s6-rc.d/init-appdaemon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: V2G Liberty
# Configures V2G Liberty
# ==============================================================================

# Creates initial AppDaemon configuration in case it is non-existing
if ! bashio::fs.file_exists '/config/appdaemon.yaml'; then
    cp -R /root/appdaemon/* /config/ \
        || bashio::exit.nok 'Failed to create initial AppDaemon configuration'
fi

if ! bashio::fs.file_exists '/config/apps/v2g-liberty'; then
    mkdir -p /config/apps/v2g-liberty
fi
cp /root/v2g-liberty/*.{py,yaml} /config/apps/v2g-liberty \
    || bashio::exit.nok 'Failed to copy V2G-Liberty app'

if ! bashio::fs.file_exists '/homeassistant/packages/v2g-liberty'; then
    mkdir -p /homeassistant/packages/v2g-liberty
fi
cp /root/v2g-liberty/package/* /homeassistant/packages/v2g-liberty \
    || bashio::exit.nok 'Failed to copy V2G-Liberty package'
